# pacman/archdots.service
[Unit]
Description=Backup package list and refresh arch/dotfiles repos

[Service]
ExecStart=/home/dylan/arch/pacman/archdots.sh
